# hogwarts-site-service/supervisor.conf
[program:hogwarts-site]
command=nodejs .
directory=/home/hogwarts_www/site
user=hogwarts_www
